# bluer_ai/assets/rpi4/bluer_ai.service
[Unit]
Description=Bluer AI Autostart (sudo -E equivalent)
After=graphical.target
Wants=graphical.target

[Service]
Type=simple

# Equivalent to sudo (run as root)
User=root

# Equivalent to sudo -E (preserve environment)
Environment=HOME=/home/pi
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=DISPLAY=:0

# Optional — preserve PATH like sudo -E does
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

WorkingDirectory=/home/pi

# Your command
ExecStart=/home/pi/git/bluer-ai/bluer_ai/.abcli/bluer_ai.sh "~terraform,where=autostart" bluer_ai session start

# Do NOT restart on failure
Restart=no

[Install]
WantedBy=graphical.target
